# healthcheck.service
[Unit]
Description=Healthcheck FastAPI Application Service
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/apps/healthcheck
Type=forking
ExecStart=/home/ubuntu/apps/healthcheck/start_app
ExecStop=/home/ubuntu/apps/healthcheck/stop_app
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target